# supervisor/conf.d/wunderground-slash-command.conf
[program:slashcommands]
command=/srv/mattermost/venv/bin/python3 /srv/mattermost/wunderground-slash-command.py
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/slashcommands
